# tools/wifi-manager/wifi-manager.service
[Unit]
Description=Robust Wi-Fi Manager for Drone Show
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi-manager.sh
WorkingDirectory=/usr/local/bin
User=root
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifi-manager

[Install]
WantedBy=multi-user.target
