# tools/led_indicator/led_indicator.service
[Unit]
Description=LED Indicator for Drone System Boot
# Run early in boot to indicate system is starting
After=sysinit.target
Before=basic.target

[Service]
Type=oneshot
RemainAfterExit=no
# Use --state for semantic state (BOOT_STARTED = red)
ExecStart=__MDS_INSTALL_DIR__/venv/bin/python __MDS_INSTALL_DIR__/led_indicator.py --state BOOT_STARTED
WorkingDirectory=__MDS_INSTALL_DIR__
User=__MDS_USER__
Group=__MDS_USER__

# GPIO access for LED control
CapabilityBoundingSet=CAP_SYS_RAWIO
AmbientCapabilities=CAP_SYS_RAWIO

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-indicator

[Install]
WantedBy=sysinit.target
