# tools/coordinator.service
[Unit]
Description=MDS Drone Swarm Coordinator
After=network-online.target git_sync_mds.service mavlink-router.service
Wants=network-online.target mavlink-router.service git_sync_mds.service
# Note: Using Wants instead of Requires for graceful degradation
# Coordinator can run even if git_sync or mavlink-router aren't available
# Restart limits (must be in [Unit] for systemd 249+)
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
ExecStart=__MDS_INSTALL_DIR__/venv/bin/python __MDS_INSTALL_DIR__/coordinator.py
WorkingDirectory=__MDS_INSTALL_DIR__
User=__MDS_USER__
Group=__MDS_USER__
Environment="HOME=__MDS_HOME__"
Environment="PATH=__MDS_INSTALL_DIR__/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Restart policy
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

# Resource limits for enterprise deployment
LimitNOFILE=65536
MemoryMax=512M

# Capabilities - reduced to minimum required
# CAP_NET_ADMIN: Network configuration for MAVLink
# CAP_NET_BIND_SERVICE: Bind to privileged ports if needed
# CAP_SYS_RAWIO: GPIO/LED hardware access
# Removed: CAP_SYS_TIME, CAP_SYS_BOOT (use NTP and systemd instead)
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_SYS_RAWIO
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_SYS_RAWIO

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=coordinator

[Install]
WantedBy=multi-user.target
